# ipfs/init.d/ipfs-daemon
#!/sbin/runscript

depend() {
  need net
}

start() {
	ebegin "Starting IPFS daemon"
	start-stop-daemon --start --background \
		--pidfile /var/run/ipfs.pid --make-pidfile \
		--exec /usr/bin/ipfs \
		--user ipfs \
		-e HOME=/pub/dldr/ipfs -- daemon
	eend $?
}

stop() {
	start-stop-daemon --stop \
		--pidfile /var/run/ipfs.pid \
		--user ipfs \
		--name ipfs
	eend $?
}
